# .config/autostart/ssh-identities-add.desktop
[Desktop Entry]
Type=Application
Name=ssh-identities-add
Exec=ssh-identities-add
